# client-system/systemd/kpanel-client.service
[Unit]
Description=KPanel Client Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/kpanel-client
Environment=KPANEL_DEVICE_ID=${HOSTNAME}
Environment=KPANEL_STATE_PATH=/var/lib/kpanel-client/device-state.json
ExecStart=/usr/local/bin/kpanel-client-launcher.sh
StateDirectory=kpanel-client
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
