# client-system/systemd/kpanel-updater.path
[Unit]
Description=Trigger KPanel updater when update-now request is queued

[Path]
PathExists=/tmp/kpanel-update-now.request

[Install]
WantedBy=multi-user.target
