# client-system/systemd/kpanel-client.service
[Unit]
Description=KPanel Client Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/kpanel-client
Environment=KPANEL_DEVICE_ID=${HOSTNAME}
Environment=KPANEL_STATE_PATH=/var/lib/kpanel-client/device-state.json
Environment=KPANEL_WIFI_ENABLED=true
Environment=KPANEL_WIFI_UI_ENABLED=true
Environment=KPANEL_WIFI_CONNECT_TIMEOUT_SEC=25
Environment=KPANEL_HOTSPOT_FALLBACK_ENABLED=true
Environment=KPANEL_HOTSPOT_SSID=KPanel-Recovery
Environment=KPANEL_HOTSPOT_PASSWORD=kpanel12345
ExecStart=/usr/local/bin/kpanel-client-launcher.sh
StateDirectory=kpanel-client
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
