# client-system/image/pi-gen/stage-kpanel/00-files/etc/systemd/system/kpanel-pi-self-heal.service
[Unit]
Description=Reassert KPanel Pi appliance defaults
DefaultDependencies=no
After=local-fs.target systemd-user-sessions.service
Before=getty@tty1.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/kpanel-pi-self-heal
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
