# client-system/image/pi-gen/stage-kpanel/00-files/home/pi/.config/autostart/kpanel-client.desktop
[Desktop Entry]
Type=Application
Name=KPanel Client
Comment=Start KPanel kiosk client at login
Exec=/usr/local/bin/kpanel-client-launcher.sh
X-GNOME-Autostart-enabled=true
NoDisplay=false
Terminal=false
